# fs-root/etc/services.d/wireguard/finish
#!/bin/bash

[[ -f /etc/wireguard/wg0.conf ]] && {
	echo "[$(date -Iseconds)] Connection summary:"
	wg show wg0 2>/dev/null
	wg-quick down wg0 2>/dev/null
	rm -f /etc/wireguard/wg0.conf
}
